# crontab
*/15 * * * * cd /app && go run schedule/stat_record.go
